# .config/systemd/user/beet-web.service
[Unit]
Description=Beets Web UI

[Service]
ExecStart=/usr/bin/beet web

[Install]
WantedBy=default.target
