# .config/systemd/user/check-battery.timer
[Unit]
Description=Check battery level

[Timer]
OnCalendar=*:0/2
AccuracySec=60

[Install]
WantedBy=timers.target
