# .config/systemd/user/check-battery.timer
[Unit]
Description=Check battery level

[Timer]
OnActiveSec=120
AccuracySec=60

[Install]
WantedBy=timers.target
