# .config/systemd/user/emacs-integration.service
[Unit]
Description=Emacs remote integration service
Requires=emacs-integration.socket

[Service]
ExecStart=%h/go/bin/emacs-integration
